# auto-remove-recheck-key.cron
0 0 1 * *     opendnssec  date >> /var/log/auto-remove-recheck-key.log && /PATH/TO/SCRIPT/auto-remove-recheck-key >> /var/log/auto-remove-recheck-key.log
